# deploy/systemd/llamafleet.service
[Unit]
Description=LlamaFleet (single service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llamafleet
Group=llamafleet
WorkingDirectory=/opt/llamafleet
EnvironmentFile=-/etc/llamafleet/llamafleet.env
ExecStart=/usr/bin/node /opt/llamafleet/scripts/run-native.mjs
Restart=always
RestartSec=2
TimeoutStopSec=60
KillMode=control-group
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/llamafleet /opt/llamafleet/data

[Install]
WantedBy=multi-user.target
